# crontab/script.sh
#!/bin/bash

curl https://tune-bot.ru/abstraction/